# docker/systemd/samba-admin-deploy@.service
[Unit]
Description=Samba Admin Docker Compose deploy service (%i)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot

# Absolute project directory used by docker compose.
# The install script replaces __PROJECT_DIR__ automatically.
WorkingDirectory=__PROJECT_DIR__

# You can override the compose binary if needed.
# Example: /usr/local/bin/docker
Environment="DOCKER_BIN=/usr/bin/docker"
Environment="GIT_BIN=/usr/bin/git"

# Deploy commands by instance name (rebuild image and recreate container):
# - samba-admin-deploy@backend => deploy backend only
# - samba-admin-deploy@frontend => deploy frontend only (without dependencies)
# - samba-admin-deploy@all => deploy all services in compose
ExecStart=/bin/bash -lc '\
  unset BACKEND_PORT FRONTEND_PORT BACKEND_DATA_PATH BACKEND_ENV_FILE \
        BACKEND_CPUS BACKEND_MEM_LIMIT BACKEND_MEM_RESERVATION \
        FRONTEND_CPUS FRONTEND_MEM_LIMIT FRONTEND_MEM_RESERVATION \
        COMPOSE_FILE COMPOSE_PROFILES COMPOSE_PROJECT_NAME; \
  compose_env_file=\"$PWD/docker/compose.env\"; \
  env_opt=\"\"; [ -f \"$compose_env_file\" ] && env_opt=\"--env-file $compose_env_file\"; \
  ${GIT_BIN} pull --ff-only && \
  case "%i" in \
    backend) ${DOCKER_BIN} compose ${env_opt} up -d --build --force-recreate backend ;; \
    frontend) ${DOCKER_BIN} compose ${env_opt} up -d --build --force-recreate --no-deps frontend ;; \
    all) ${DOCKER_BIN} compose ${env_opt} up -d --build --force-recreate ;; \
    *) echo "Invalid instance: %i (use backend|frontend|all)"; exit 1 ;; \
  esac'

# Allow enough time for image build and container recreation.
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
