# docker/systemd/samba-admin@.service
[Unit]
Description=Samba Admin Docker Compose service (%i)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Absolute project directory used by docker compose.
# The install script replaces __PROJECT_DIR__ automatically.
WorkingDirectory=__PROJECT_DIR__

# You can override the compose binary if needed.
# Example: /usr/local/bin/docker
Environment="DOCKER_BIN=/usr/bin/docker"

# Start commands by instance name:
# - samba-admin@backend => start only backend
# - samba-admin@frontend => start only frontend (without dependencies)
# - samba-admin@all => start all services in compose
ExecStart=/bin/bash -lc '\
  case "%i" in \
    backend) ${DOCKER_BIN} compose up -d backend ;; \
    frontend) ${DOCKER_BIN} compose up -d --no-deps frontend ;; \
    all) ${DOCKER_BIN} compose up -d ;; \
    *) echo "Invalid instance: %i (use backend|frontend|all)"; exit 1 ;; \
  esac'

# Stop commands by instance name.
# For "all" we use "down" to stop and remove the compose stack.
ExecStop=/bin/bash -lc '\
  case "%i" in \
    backend) ${DOCKER_BIN} compose stop backend ;; \
    frontend) ${DOCKER_BIN} compose stop frontend ;; \
    all) ${DOCKER_BIN} compose down ;; \
    *) echo "Invalid instance: %i (use backend|frontend|all)"; exit 1 ;; \
  esac'

# Allow enough time for graceful stop/start actions.
TimeoutStartSec=300
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
