# docker/systemd/samba-admin@.service
[Unit]
Description=Samba Admin Docker Compose service (%i)
Requires=docker.service
After=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Absolute project directory used by docker compose.
# The install script replaces __PROJECT_DIR__ automatically.
WorkingDirectory=__PROJECT_DIR__

# You can override the compose binary if needed.
# Example: /usr/local/bin/docker
Environment="DOCKER_BIN=/usr/bin/docker"

# Start commands by instance name:
# - samba-admin@backend => start only backend
# - samba-admin@frontend => start only frontend (without dependencies)
# - samba-admin@all => start all services in compose
ExecStart=/bin/bash -c '\
  compose_env_file=\"__PROJECT_DIR__/docker/compose.env\"; \
  if [ ! -f \"$compose_env_file\" ]; then \
    echo \"Missing compose env file: $compose_env_file\"; \
    exit 1; \
  fi; \
  case "%i" in \
    backend) ${DOCKER_BIN} compose --project-directory \"__PROJECT_DIR__\" -f \"__PROJECT_DIR__/docker-compose.yml\" --env-file \"$compose_env_file\" up -d backend ;; \
    frontend) ${DOCKER_BIN} compose --project-directory \"__PROJECT_DIR__\" -f \"__PROJECT_DIR__/docker-compose.yml\" --env-file \"$compose_env_file\" up -d --no-deps frontend ;; \
    all) ${DOCKER_BIN} compose --project-directory \"__PROJECT_DIR__\" -f \"__PROJECT_DIR__/docker-compose.yml\" --env-file \"$compose_env_file\" up -d ;; \
    *) echo "Invalid instance: %i (use backend|frontend|all)"; exit 1 ;; \
  esac'

# Stop commands by instance name.
# For "all" we use "down" to stop and remove the compose stack.
ExecStop=/bin/bash -c '\
  compose_env_file=\"__PROJECT_DIR__/docker/compose.env\"; \
  if [ ! -f \"$compose_env_file\" ]; then \
    echo \"Missing compose env file: $compose_env_file\"; \
    exit 1; \
  fi; \
  case "%i" in \
    backend) ${DOCKER_BIN} compose --project-directory \"__PROJECT_DIR__\" -f \"__PROJECT_DIR__/docker-compose.yml\" --env-file \"$compose_env_file\" stop backend ;; \
    frontend) ${DOCKER_BIN} compose --project-directory \"__PROJECT_DIR__\" -f \"__PROJECT_DIR__/docker-compose.yml\" --env-file \"$compose_env_file\" stop frontend ;; \
    all) ${DOCKER_BIN} compose --project-directory \"__PROJECT_DIR__\" -f \"__PROJECT_DIR__/docker-compose.yml\" --env-file \"$compose_env_file\" down ;; \
    *) echo "Invalid instance: %i (use backend|frontend|all)"; exit 1 ;; \
  esac'

# Allow enough time for graceful stop/start actions.
TimeoutStartSec=300
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
